# ddns-agent.service
[Unit]
Description=DDNS Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ddns-agent
Environment=AGENT_CONFIG_PATH=/etc/ddns-agent/config.enc.json
Environment=AGENT_DB_PATH=/var/lib/ddns-agent/agent.db
EnvironmentFile=/etc/ddns-agent/agent.env
WorkingDirectory=/var/lib/ddns-agent
ExecStart=/bin/bash -lc 'exec "${DDNS_WORKDIR:-/var/lib/ddns-agent}/start-agent.sh"'
Restart=on-failure
ProtectSystem=strict
ProtectHome=yes
BindReadOnlyPaths=/etc/ddns-agent/config.enc.json
BindPaths=/var/lib/ddns-agent
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
